# PLPProxyClient/etc/frp/frpc.service
[Unit]
Description=frp Client
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/frpc -c /etc/frp/frpc.toml
Restart=always
RestartSec=5
User=nobody

[Install]
WantedBy=multi-user.target
